# rpm/longview_py.cron
*/5 * * * * root /opt/longview_py/longview.py > /tmp/longview_py.lastlog 2>&1
